# crontab
0 * * * * cd /root/sentry && npm run dev

